# etc/supervisor.d/apache.ini
[program:apache]
command = /usr/sbin/httpd -DFOREGROUND
user = root
autostart = false