# etc/supervisor.d/gitlab-ci-runner.ini
[program:gitlab-ci-runner]
command=/usr/bin/gitlab-ci-multi-runner run
user=root
autorestart=true