# etc/supervisor.d/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
pidfile = /var/run/nginx.pid
user = root
autostart = false