# etc/supervisor.d/php5-fpm.ini
[program:php5-fpm]
command = /usr/bin/php-fpm -FR
pidfile = /var/run/php5-fpm.pid
user = root
autostart = true