# etc/supervisor.d/php5-fpm.ini
[program:php5-fpm]
command = /usr/bin/php-fpm -FR
user = root
autostart = true