# etc/supervisor.d/php7-fpm.ini
[program:php7-fpm]
command = /usr/sbin/php-fpm7 -FR
pidfile = /var/run/php7-fpm.pid
user = root
autostart = true