# etc/supervisor.d/php7-fpm.ini
[program:php7-fpm]
command = /usr/sbin/php-fpm7 -FR
user = root
autostart = true