# etc/supervisor.d/setup-xdebug.ini
[program:setup-xdebug]
command = /usr/local/bin/setup-xdebug -n -r
startsecs = 0
autorestart = false
startretries = 1