# tmp/etc/supervisor.d/mysql.ini
[program:mysql]
command = /usr/bin/pidproxy /run/mysqld/mysqld.pid /usr/bin/mysqld
pidfile = /run/mysqld/mysqld.pid
user = mysql
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0