# tmp/etc/supervisor.d/postfix.ini
[program:postfix]
command = /usr/sbin/postfix -c /etc/postfix start
pidfile = /var/run/postfix.pid
user = root
directory = /etc/postfix
startsecs = 0
autorestart	= false
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0